# services/lpmc_smp.service
[Service]
ExecStart=bash /home/ubuntu/LPMC/scripts/start.sh smp
Restart=on-failure
WorkingDirectory=/home/ubuntu/LPMC
