# services/lpmc_waterfall.service
[Service]
ExecStart=bash /home/ubuntu/LPMC/scripts/start.sh waterfall
Restart=on-failure
WorkingDirectory=/home/ubuntu/LPMC
