# services/lpmc_lobby.service
[Service]
ExecStart=/home/ubuntu/LPMC/scripts/start.sh lobby
Restart=on-failure
WorkingDirectory=/home/ubuntu/LPMC
